# powerdown@.service
[Service]
Type=oneshot
ExecStart=/usr/bin/pm-powersave %I
